# supervisord.conf
[supervisord]
nodaemon=true

[program:rabbitmq]
command=/bin/sh /opt/rabbitmq.sh
user=root
autostart=true
autorestart=true
environment=
  HOME="/var/lib/rabbitmq/",
  RABBITMQ_PID_FILE="/var/run/rabbitmq/pid"

[program:server]
command=python /app/TweetFeedPlus.py

[program:google-fluentd]
command=/bin/bash /install-logging-agent.sh